# services/gollum/run
#!/bin/bash

cd /data

/sbin/setuser abc gollum 2>&1
